# profile.d/progress.sh
DLC=/usr/dlc
PATH=$PATH:$DLC/bin
export DLC PATH
