# profile.d/progress.sh
# Test same commit on 2 locations
DLC=${DLC:-/usr/dlc}
PATH=$PATH:$DLC/bin
export DLC PATH
